# deploy/reinventa-worker.service
[Unit]
Description=Reinventa Laravel queue worker
After=network.target mysql.service

[Service]
Type=simple
User=www-data
Group=www-data
Restart=always
RestartSec=5
WorkingDirectory=/var/www/reinventa
ExecStart=/usr/bin/php8.3 artisan queue:work --tries=5 --backoff=30 --max-time=3600 --sleep=3

[Install]
WantedBy=multi-user.target
